# Inception/srcs/mariadb/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld
autostart=true 
autorestart=true

[program:telegraf]
command = telegraf
autostart=true 
autorestart=true